# src/supervisord.conf
[supervisord]
nodaemon=false
logfile=Logs/supervisord.log
pidfile=Logs/supervisord.pid
childlogdir=Logs

[unix_http_server]
file=Logs/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix://Logs/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; Gmail Watcher
[program:gmail-watcher]
command=python src/watcher_gmail.py
directory=%(here)s
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=10
stdout_logfile=Logs/supervisor_gmail_out.log
stderr_logfile=Logs/supervisor_gmail_err.log
environment=PYTHONUNBUFFERED="1"

; WhatsApp Watcher
[program:whatsapp-watcher]
command=python src/watcher_whatsapp.py
directory=%(here)s
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=10
stdout_logfile=Logs/supervisor_whatsapp_out.log
stderr_logfile=Logs/supervisor_whatsapp_err.log
environment=PYTHONUNBUFFERED="1"

; Orchestrator
[program:orchestrator]
command=python orchestrator.py --watch --interval 300
directory=%(here)s
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=10
stdout_logfile=Logs/supervisor_orchestrator_out.log
stderr_logfile=Logs/supervisor_orchestrator_err.log
environment=PYTHONUNBUFFERED="1"

; Health Monitor
[program:health-monitor]
command=python Scripts/health_monitor.py
directory=%(here)s
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=10
stdout_logfile=Logs/supervisor_health_out.log
stderr_logfile=Logs/supervisor_health_err.log
environment=PYTHONUNBUFFERED="1"

; Group all processes
[group:ai-employee]
programs=gmail-watcher,whatsapp-watcher,orchestrator,health-monitor
priority=999
